# com.openwebui.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openwebui.service</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>set -a; source /opt/Open-WebUI-Service/.env; set +a; exec /opt/Open-WebUI-Service/.venv/bin/open-webui serve</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/Open-WebUI-Service</string>

    <key>UserName</key>
    <string>root</string>

    <key>GroupName</key>
    <string>wheel</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/opt/Open-WebUI-Service/openwebui.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/opt/Open-WebUI-Service/openwebui.err.log</string>
</dict>
</plist>
